# mfsbsdonly/etc/rc.local.sample
# sample rc.local 
# add your code below